# aurore/supervisor/supervisord-dev.conf
[supervisord]
nodaemon=true
autostart=true
autorestart=true
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=3MB
logfile_backups=1
loglevel=debug
pidfile=/tmp/supervisord.pid

[program:aurore]
command=gunicorn -c /usr/src/aurore/gunicorn/gunicorn-dev.conf.py --chdir /usr/src/aurore/src aurore.wsgi
autostart=true
autorestart=true
directory= /usr/src/aurore/src

[program:interrupt]
command=python /usr/src/aurore/src/interrupt.py
autostart=true
autorestart=true
directory= /usr/src/aurore/src

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=ludock:ludock
username=ludock
password=%(ENV_SV_PWD)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf